# iptv-player.service
[Unit]
Description=IPTV Smart Player for Granny TV
After=network.target graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
User=jeremy
Group=jeremy
WorkingDirectory=/home/jeremy/pi
ExecStart=/usr/bin/python3 /home/jeremy/pi/iptv_smart_player.py
Restart=always
RestartSec=10
Environment=DISPLAY=:0
Environment=HOME=/home/jeremy

# Logging
StandardOutput=append:/home/jeremy/pi/iptv_service.log
StandardError=append:/home/jeremy/pi/iptv_service.log

[Install]
WantedBy=graphical-session.target